# tests/systemd/test3.service
[Unit]
Description=test3

[Service]
Type=oneshot
Restart=no
ExecStart=/usr/libexec/test3.sh

[Install]
WantedBy=multi-user.target
