# tests/systemd/test1.service
[Unit]
Description=test1

[Service]
Type=oneshot
Restart=no
ExecStart=/usr/libexec/test1.sh

[Install]
WantedBy=multi-user.target
